# hdo-app.conf
description "HDO-Transcript Search"

start on startup
stop on shutdown

respawn
respawn limit 10 90

exec start-stop-daemon --start -c your-username --exec npm start /path/to/hdo/hdo-app.sh
